# system/update.timer
[Unit]
Description=Apt Update

[Timer]
Unit=update.service
OnUnitActiveSec=8min
OnBootSec=8min

[Install]
WantedBy=timers.target
