# system/update.service
[Unit]
Description=Apt Update

[Service]
Type=simple
User=root
ExecStart=/usr/bin/apt-get update
