# system/upgrade.service
[Unit]
Description=Apt Upgrade

[Service]
Type=simple
User=root
ExecStart=/usr/bin/apt-get upgrade --yes
