# system/upgrade.timer
[Unit]
Description=Apt Upgrade

[Timer]
Unit=upgrade.service
OnUnitActiveSec=12min
OnBootSec=12min

[Install]
WantedBy=timers.target
